# bots/bosun/systemd.service
[Unit]
Description=Pacto bot bosun
After=network.target

[Service]
Type=simple
User=pacto-bot
Group=pacto-bot
WorkingDirectory=/opt/bosun
ExecStart=/usr/local/bin/python3 bosun.py
Restart=on-failure
Environment="PACTO_TRANSPORT=unix"
Environment="PACTO_SOCKET_PATH=/run/pacto-bot-api.sock"

[Install]
WantedBy=multi-user.target
